# src/cephadm/cephadmlib/templates/dropin.service.j2
# generated by cephadm
[Unit]
{%- for volume in volume_ids %}
Wants={{ volume.volume_service_name }}
{%- endfor %}
{%- if enable_init_containers %}
Wants={{ identity.init_service_name }}
{%- endif %}
{%- for sidecar in sidecar_ids %}
Wants={{ sidecar.sidecar_service_name }}
{%- endfor %}

