# src/cephadm/cephadmlib/templates/volume.service.j2
# generated by cephadm
[Unit]
Description=Volume for Ceph Service %i
After=network-online.target local-fs.target time-sync.target
Wants=network-online.target local-fs.target time-sync.target
{%- if has_docker_engine %}
After=docker.service
Wants=docker.service
{%- endif %}

Before={{ primary.service_name }}
{%- if enable_init_containers %}
Before={{ primary.init_service_name }}
{%- endif %}


[Service]
EnvironmentFile=-/etc/environment
ExecStartPre=-{{ volume.rm_command(ctx.container_engine.path) | map('shellquote') | join(' ') }}
ExecStart={{ volume.create_command(ctx.container_engine.path) | map('shellquote') | join(' ') }}
ExecStop={{ volume.rm_command(ctx.container_engine.path) | map('shellquote') | join(' ') }}
TimeoutStopSec=120
StartLimitInterval=30min
StartLimitBurst=5
Type=oneshot
RemainAfterExit=yes
